# nvidia-unbind-check.service
[Unit]
Description=Check and unbind NVIDIA from VFIO if no VMs running
After=libvirtd.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/check-qemu-unbind.sh
